# ansible/todoapp.service
[Unit]
Description=My To-Do App

[Service]
User=ec2-user
WorkingDirectory=/opt/todoapp
ExecStart=/opt/todoapp/poetry run flask run --host 0.0.0.0